# packaging/debian/postinst
#!/bin/sh
# Post-install: create the root-only config dir and start the enforcement
# service so blocking is active immediately and on every boot.
set -e

CFG=/etc/appblocker

case "$1" in
    configure)
        mkdir -p "$CFG"
        chmod 711 "$CFG"   # children can't list/read the password or blocklist

        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload || true
            # Always (re-)enable for boot and (re)start the new binary. Explicit
            # enable + restart also heals any machine left disabled by an older
            # buggy prerm, and guarantees the freshly-installed binary is running.
            systemctl enable appblocker.service || true
            systemctl restart appblocker.service || true
        else
            echo "systemd not detected; start the daemon manually with:" >&2
            echo "  /usr/bin/appblocker --daemon" >&2
        fi
        ;;
esac

exit 0
