# packaging/debian/prerm
#!/bin/sh
# Pre-removal: stop and disable the enforcement service.
set -e

case "$1" in
    remove|deconfigure)
        if [ -d /run/systemd/system ]; then
            systemctl disable --now appblocker.service || true
        fi
        # Remove website blocks from /etc/hosts and browser lockdown policies.
        [ -x /usr/bin/appblocker ] && /usr/bin/appblocker --web-clear || true
        [ -x /usr/bin/appblocker ] && /usr/bin/appblocker --lockdown-clear || true
        ;;
    upgrade)
        if [ -d /run/systemd/system ]; then
            systemctl disable --now appblocker.service || true
        fi
        ;;
esac

exit 0
